# deployment/crontab
 * * * * * root ls camellia-stats-scraper/deployment/run-scraper-in-docker > /tmp/hello
 * * * * * root ls > /tmp/hello2
 