# deployment/crontab
 */30 * * * * camellia-stats/deployment/run-scraper-in-docker prod > logs/stdout 2>logs/stderr
 
