# deployment/crontab
 * * * * * root ls camellia-stats-scraper/deployment/run-scraper-in-docker
 * * * * * root ls 
 