# deployment/crontab
 * * * * * root camellia-stats-scraper/deployment/run-scraper-in-docker
 